# home/dot_config/systemd/user/ollama.service
[Unit]
Description=Ollama
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/ollama serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
